# resources/pulseaudio/pulseaudio.socket
[Socket]
ListenStream=172.20.0.1:4713

[Install]
WantedBy=sockets.target
